# deploy/bg-remover.service
[Unit]
Description=BG Remover API Server
After=network.target

[Service]
Type=simple
EnvironmentFile=/opt/bg-remover/.env
ExecStart=/opt/bg-remover/bg-remover-server
WorkingDirectory=/opt/bg-remover
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
